# Docker/supervisord.conf
[supervisord]
logfile=/app/host/logs/supervisord/supervisord.log ; supervisord log file
logfile_maxbytes=50MB                         ; maximum size of logfile before rotation
logfile_backups=7                             ; number of backed up logfiles
loglevel=error                                ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid              ; pidfile location
nodaemon=false                                ; run supervisord as a daemon
minfds=1024                                   ; number of startup file descriptors
minprocs=200                                  ; number of process descriptors
user=root                                     ; default user

[supervisorctl]

[program:CosmWebex]
command=/bin/bash /app/app-run.sh
autostart=true
autorestart=true
stderr_logfile=/app/host/logs/supervisord/CosmWebex-supervisor.err.log
stdout_logfile=/app/host/logs/supervisord/CosmWebex-supervisor.out.log
